# pi-client/starfire-dev.service
[Unit]
Description=星火智造云打印 - Vite Dev Server (热更新)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/StarfirePrintHub/frontend
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/starfire-dev.log
StandardError=append:/var/log/starfire-dev.log

# 环境变量
Environment=NODE_ENV=development
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
